# systemd/cloudflared-tunnel.service
# Install as: /etc/systemd/system/cloudflared-TUNNEL.service
# Replace TUNNEL with your tunnel name and YOUR_USER with your LAB_USER, then:
#   sudo systemctl daemon-reload && sudo systemctl enable --now cloudflared-TUNNEL
# (bootstrap.sh does this substitution for you from lab.conf.)
[Unit]
Description=Cloudflare Tunnel (TUNNEL)
After=network-online.target caddy.service
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
ExecStart=/usr/local/bin/cloudflared --no-autoupdate --config /home/YOUR_USER/lab/cloudflared/config.yml tunnel run TUNNEL
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
