# systemd/flame-status.timer
# Install as: /etc/systemd/system/flame-status.timer  then:
#   sudo systemctl daemon-reload && sudo systemctl enable --now flame-status.timer
[Unit]
Description=Probe apps and refresh dashboard status dots every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=5s

[Install]
WantedBy=timers.target
